# deploy/sms-gate.service
[Unit]
Description=SMS Gate API
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5
OnFailure=sms-gate-notify@%n.service

[Service]
Type=simple
User=smsgate
Group=dialout
WorkingDirectory=/opt/sms-gate
ExecStart=/opt/sms-gate/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 30080
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/sms-gate/.env

[Install]
WantedBy=multi-user.target
